# crontab
56 21 * * 0 /home/mzkwcim/sender/run_python_script.sh > /home/mzkwcim/sender/cronjob.log 2>&1
